# services/sampler.service
[Unit]
Description=Network Multimeter Sampler

[Service]
ExecStart=/usr/bin/sampler.sh eth0 /tmp/capture

[Install]
WantedBy=multi-user.target

